# aws/systemd/proofport-ai.service
[Unit]
Description=proofport-ai Blue-Green Container Manager
Documentation=https://github.com/zkproofport/proofport-ai
After=docker.service network-online.target proofport-ai-redis.service
Wants=network-online.target
Requires=docker.service proofport-ai-redis.service

[Service]
Type=oneshot
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal
SyslogIdentifier=proofport-ai

# Environment file — contains AI_IMAGE, DEPLOY_ENV, and all app config.
# Operator must populate /opt/proofport-ai/.env before enabling this service.
EnvironmentFile=/opt/proofport-ai/.env

# On boot: start the last-active blue-green slot container.
# Reads /opt/proofport-ai/active-slot to determine which slot (blue or green)
# was last deployed, then starts the corresponding container on the correct ports.
ExecStart=/opt/proofport-ai/boot-active-slot.sh

# On stop: gracefully stop the active container.
ExecStop=/opt/proofport-ai/stop-active-slot.sh
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
